# files/etc/crontabs/root
0 4 * * * /etc/coremark.sh
0 2 * * * /usr/share/UnblockNeteaseMusic/update_core.sh 2>&1
